# level2/level2.socket
[Unit]
Description=level2 puzzlebox challenge

[Socket]
ListenStream=0.0.0.0:1338
Accept=true

[Install]
WantedBy=sockets.target
